# private_Library/LaunchAgents/ai.factory.droid.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>ai.factory.droid.daemon</string>

	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/droid</string>
		<string>daemon</string>
		<string>--remote-access</string>
	</array>

	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/jonah.snider</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>

	<key>KeepAlive</key>
	<true/>

	<key>RunAtLoad</key>
	<true/>

	<key>StandardErrorPath</key>
	<string>/Users/jonah.snider/Library/Logs/factory-droid-daemon.log</string>

	<key>StandardOutPath</key>
	<string>/Users/jonah.snider/Library/Logs/factory-droid-daemon.log</string>

	<key>WorkingDirectory</key>
	<string>/Users/jonah.snider</string>
</dict>
</plist>
